# overlays/etc/systemd/system/phosh.service
[Unit]
Description=Phosh Mobile Shell
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
User=droidian
Group=droidian
ExecStart=/usr/bin/phosh-session
Restart=always
RestartSec=5
Environment=XDG_SESSION_TYPE=wayland
Environment=XDG_SESSION_DESKTOP=phosh
Environment=XDG_CURRENT_DESKTOP=phosh

[Install]
WantedBy=graphical.target